# deploy/regime-trader.service
[Unit]
Description=regime_trader — one daily decision cycle (HMM/trend-core paper bot)
Documentation=file:///opt/regime_trader/deploy/README.md
# The bot's first act is to reach api.alpaca.markets.
After=network-online.target
Wants=network-online.target

[Service]
# One decision cycle per ticker, then exit — driven by regime-trader.timer.
# A daily-bar system only decides once a day, so the process runs for ~15s
# and the host idles the rest of the day. This is NOT enabled directly;
# enable regime-trader.timer instead.
Type=oneshot
User=regime
Group=regime
WorkingDirectory=/opt/regime_trader
EnvironmentFile=/opt/regime_trader/.env
ExecStart=/opt/regime_trader/.venv/bin/python -u main.py --once

# EXIT_HALTED (3): the -20% drawdown breaker fired. --once checks the lock
# BEFORE any broker call and exits 3, so the next daily fire is a harmless
# no-op until a human clears it. SuccessExitStatus keeps the unit from being
# spammed as "failed" every day while halted (it is a known, surfaced state;
# the healthcheck reports the lock separately).
SuccessExitStatus=0 3

# ── Sandboxing / hardening ─────────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/regime_trader/logs /opt/regime_trader/data_cache
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

# No [Install] section: this unit is triggered by the timer, not enabled
# on its own.
