# deploy/regime-trader.timer
[Unit]
Description=Fire regime_trader once per trading day (shortly after the US open)
Documentation=file:///opt/regime_trader/deploy/README.md

[Timer]
# 09:35 America/New_York, Mon-Fri. The explicit timezone makes systemd
# handle EST/EDT itself, so this stays 5 minutes after the 09:30 open all
# year (a naive UTC time would drift an hour across DST). US market holidays
# still fire but are harmless no-ops: the market clock reads closed and the
# latest daily bar is stale, so no order is placed.
#
# Why just after the open (not the close): the existing loop submits an
# immediate market order on the latest bar, so it decides on today's fresh
# bar and fills right away. (A close-time run submitting market-on-open
# orders for the next day would map the backtest's next-open fills even more
# exactly — a possible future refinement, see README.)
OnCalendar=Mon-Fri 09:35 America/New_York

# If the host was asleep/off at fire time (e.g. a reboot window), run as soon
# as it is back rather than skipping the day. A slow trend decision made a
# few hours late is far better than a missed rebalance.
Persistent=true

# Small randomised delay so many hosts don't hit the API on the same second.
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
